# shared_files/aryaos/systemd/aryaos-neighbord.service
[Unit]
Description=AryaOS CoT neighbor discovery cache
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/aryaos/aryaos-config.txt
ExecStart=/usr/local/sbin/aryaos-neighbord
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=multi-user.target
