# shared_files/aryaos/systemd/aryaos-boot-stable.service
[Unit]
Description=AryaOS boot-stable marker — reset the crash counter once the box stays up
Documentation=https://github.com/snstac/aryaos

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/aryaos-safe-mode boot-ok
