# shared_files/aryaos/systemd/gpsd.socket.d/socket-group.conf
# Allow local services (e.g. lighttpd www-data portal CGI via gpspipe) to read gpsd JSON.
# Default gpsd.socket uses SocketMode=0600 (root only).
#
# SPDX-License-Identifier: Apache-2.0
# Copyright Sensors & Signals LLC https://www.snstac.com/

[Socket]
SocketMode=0660
SocketGroup=gpsd
