# shared_files/adsbcot/systemd/dump1090-fa.service.d/aryaos-conflicts.conf
# SPDX-License-Identifier: Apache-2.0
# Copyright Sensors & Signals LLC https://www.snstac.com/
#
# Only one 1090 MHz decoder may run; starting dump1090-fa stops readsb if it was active.

[Unit]
Conflicts=readsb.service
