# shared_files/aryaos/systemd/aryaos-bt-pan.service.d/emcon.conf
# EMCON (aryaos-radio silence): do not bring up the Bluetooth PAN when radio
# silence is active (Bluetooth is rfkill-blocked). Cleared by `aryaos-radio
# silence off`.
[Unit]
ConditionPathExists=!/etc/aryaos/emcon
