# dronehone-bridge/debian/dronehone-bridge.service
[Unit]
Description=DroneHone Bluetooth Remote ID sensor bridge
Documentation=file:/usr/share/doc/dronehone-bridge/README.md
After=network-online.target bluetooth.target docker.service dronecot.service
Wants=network-online.target bluetooth.target
ConditionPathExists=/sys/class/bluetooth

[Service]
Type=simple
ExecStart=/usr/bin/dronehone-bridge
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
