# shared_files/dhbridge/systemd/dhbridge.service.d/aryaos-bt-ready.conf
[Unit]
Requires=aryaos-bt-ready.service
After=aryaos-bt-ready.service
