# shared_files/dhbridge/systemd/dhbridge.service.d/aryaos-hostname-broadcast.conf
[Service]
ExecStart=
ExecStart=/usr/local/sbin/aryaos-dhbridge-start.sh
