# shared_files/charontak/systemd/charontak.service.d/aryaos-config.conf
[Service]
EnvironmentFile=-/etc/aryaos/aryaos-config.txt
