# shared_files/aryaos/systemd/nodered.service.d/aryaos.conf
# AryaOS override for the upstream nodered.service.
#
# The node-red linux-installer ships the unit with User=root / WorkingDirectory=/root,
# which makes Node-RED run AS ROOT out of /root/.node-red — a stock settings.js with
# NO adminAuth. The result is an unauthenticated, root-privileged admin API (anyone who
# can reach :1880 can overwrite flows = RCE as root). AryaOS installs the hardened
# config (adminAuth + the palette/flows) into /home/node-red/.node-red, so pin the
# service to the node-red user and that home directory.
[Service]
User=node-red
Group=node-red
WorkingDirectory=/home/node-red
Environment="HOME=/home/node-red"
