# shared_files/aryaos/systemd/aryaos-bt-ready.service.d/emcon.conf
# EMCON (aryaos-radio silence): do not run the Bluetooth prep when radio silence
# is active — aryaos-bt-ready unblocks Bluetooth (`rfkill unblock bluetooth`),
# which would defeat the block. Cleared by `aryaos-radio silence off`.
[Unit]
ConditionPathExists=!/etc/aryaos/emcon
