# shared_files/adsbcot/systemd/readsb.service.d/aryaos-conflicts.conf
# SPDX-License-Identifier: Apache-2.0
# Copyright Sensors & Signals LLC https://www.snstac.com/
#
# Only one 1090 MHz decoder may run; starting readsb stops dump1090-fa if it was active.

[Unit]
Conflicts=dump1090-fa.service
