# shared_files/aryaos/systemd/lighttpd.service.d/aryaos-netlink.conf
# Allow iproute2 "ip" (netlink) for portal CGI — base lighttpd unit restricts families.
# See: systemd RestrictAddressFamilies in /usr/lib/systemd/system/lighttpd.service
#
# SPDX-License-Identifier: Apache-2.0
# Copyright Sensors & Signals LLC https://www.snstac.com/

[Service]
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
