# shared_files/aryaos/systemd/comitup.service.d/emcon.conf
# EMCON (aryaos-radio silence): do not start the WiFi hotspot when radio silence
# is active. Bringing up the AP makes NetworkManager re-enable the WiFi radio,
# which would defeat the rfkill block. The flag is cleared by `aryaos-radio
# silence off`, which then starts comitup again.
[Unit]
ConditionPathExists=!/etc/aryaos/emcon
