# shared_files/adsbcot/systemd/dump978-fa.service.d/aryaos-uat-serial.conf
# SPDX-License-Identifier: Apache-2.0
# Copyright Sensors & Signals LLC https://www.snstac.com/
#
# Merge UAT RTL serial from /etc/aryaos/aryaos-config.txt into /etc/default/dump978-fa
# on each service start so operators can override without re-imaging.

[Service]
ExecStartPre=+/usr/local/libexec/aryaos/apply-dump978-serial.sh
